# python-agent/systemd/iscope-monitor.service
[Unit]
Description=iScope 360 Monitor
After=network-online.target iscope-supervisor.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/iscope-agent
EnvironmentFile=-/etc/iscope/agent.env
EnvironmentFile=-/etc/iscope-agent/agent.env
ExecStart=/opt/iscope-agent/venv/bin/python -m monitor.main
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=iscope-monitor
NoNewPrivileges=false
ProtectSystem=false

[Install]
WantedBy=multi-user.target
