# python-agent/systemd/iscope-supervisor.service
[Unit]
Description=iScope 360 Supervisor
After=network-online.target
Wants=network-online.target iscope-agent.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/iscope-agent
ExecStartPre=-/bin/bash /opt/iscope-agent/check-deps.sh
ExecStart=/bin/bash /opt/iscope-agent/supervisor_bootstrap.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=iscope-supervisor

# Environment
EnvironmentFile=-/etc/iscope/agent.env

# Security hardening
NoNewPrivileges=false
ProtectSystem=false

[Install]
WantedBy=multi-user.target
